# get.sh
# chkconfig: 3 67 22
#!/home/pi/pi_ublox-ntripsc
# description: ntrip_client
/home/pi/pi_ublox-ntrip/ntrip_client
